# data/cron.d/100-docker-cron
* * * * *   (cd /var/www/html && ./poirot workers default_worker &)
